# ubuntu/config/.config/autostart/krusader.desktop
[Desktop Entry]
Type=Application
Exec=krusader
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Krusader
Comment=
